# daemons/pulse_daemon.service
[Unit]
After=pulseaudio.service

[Service]
ExecStart=/usr/local/bin/pulse_daemon.sh

[Install]
WantedBy=default.target
